# meta-raspberry-software/recipes-tcp-server/tcp-server/files/tcp-server.service
[Unit]
Description=My custom startup program
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tcp-server
Restart=on-failure

[Install]
WantedBy=multi-user.target